# scripts/systemd/ghostshell-daemon.service
[Unit]
Description=ghostshell session recording collector (root-only central store)
Documentation=man:ghostshell(1)
After=network.target

[Service]
Type=simple
ExecStart=/usr/libexec/ghostshell-daemon
Environment=GHOSTSHELL_CONFIG=/etc/ghostshell/ghostshell.conf
# Restart on crash, but pace restarts so a daemon that dies on a bad config
# does not spin in a tight loop. RestartSec is given explicit seconds units.
Restart=on-failure
RestartSec=2s
# The daemon must run as root to own /var/lib/ghostshell and read user-local
# spool files during startup ingest.
User=root
# Hardening: it only needs its socket (/run/ghostshell-daemon.sock), the central store
# (/var/lib/ghostshell), its logfile (/var/log/ghostshell), and read+delete access to
# each user's ~/.local/share/ghostshell spool during startup ingest.
LogsDirectory=ghostshell
LogsDirectoryMode=0750
# ProtectSystem=full keeps /usr, /boot and /etc read-only; ReadWritePaths
# re-opens the dirs the daemon writes. /run is a symlink to /var/run on most
# distros and is not covered by ProtectSystem, but is listed for clarity so the
# socket can be (re)created.
ProtectSystem=full
ReadWritePaths=/var/lib/ghostshell /var/log/ghostshell /run /var/run
NoNewPrivileges=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX
PrivateDevices=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
UMask=0077
MemoryMax=512M
TasksMax=256
# --- Additional hardening (all verified safe for this daemon) -----------------
# The daemon: opens a unix socket, reads /etc/passwd + user home spools, writes
# 0600 recordings, and sets the key immutable via ioctl(FS_IOC_SETFLAGS)
# (chattr +i). None of the directives below touch those operations.
#
# No kernel-module load/unload — deny it.
ProtectKernelModules=yes
# Never reads the kernel ring buffer (/dev/kmsg, syslog(2)).
ProtectKernelLogs=yes
# Never calls personality(2); lock the execution domain.
LockPersonality=yes
# Never uses realtime scheduling.
RestrictRealtime=yes
# Never creates namespaces (nor do the aws/gsutil/rsync backup children).
RestrictNamespaces=yes
# Never sets the system clock (settimeofday/adjtimex/clock_settime).
ProtectClock=yes
# Never sets the hostname/domainname.
ProtectHostname=yes
# Syscall allow-list. @system-service is the recommended baseline for a normal
# service: it PERMITS ioctl (so chattr +i / FS_IOC_SETFLAGS still works), execve
# (so backup subprocesses still spawn), and file/socket I/O, while denying
# @mount/@reboot/@module/@swap/@obsolete/@raw-io/@privileged. Requires
# SystemCallArchitectures=native (set above) to be effective.
SystemCallFilter=@system-service
# NOTE: ProtectHome must NOT be enabled. The collector reads each user's
# ~/.local/share/ghostshell spool and DELETES the file after ingesting it into the
# central store. ProtectHome=read-only (or tmpfs) would make os.Remove of the
# spool fail silently, so stray recordings would accumulate and be re-ingested
# (duplicated) on every restart. Cross-user home access is core to this daemon.
ProtectHome=no

[Install]
WantedBy=multi-user.target
